# systemd/hyprpolkitagent.service.d/50-restart-resilience.conf
[Unit]
StartLimitIntervalSec=60
StartLimitBurst=15

[Service]
Restart=on-failure
RestartSec=2
